# deploy/webhook/kanvas-deploy-hook.service
[Unit]
Description=Kanvas deploy webhook receiver
After=network.target docker.service
Wants=docker.service

[Service]
ExecStart=/usr/bin/python3 /opt/kanvas-deploy/receiver.py
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
